# debian/postinst
#!/bin/sh
adduser --system mcctl --home /usr/lib/mcctl
systemctl enable mcctl
